# deploy/kangaroo-review-metrics.service
[Unit]
Description=Kangaroo Review metrics API
After=network.target

[Service]
Type=simple
User=docs
WorkingDirectory=/home/docs/kangaroo-review
EnvironmentFile=-/home/docs/kangaroo-review-runtime/metrics.env
Environment=KANGAROO_METRICS_HOST=127.0.0.1
Environment=KANGAROO_METRICS_PORT=18081
Environment=KANGAROO_DB_PATH=/home/docs/kangaroo-review-runtime/metrics.sqlite3
Environment=KANGAROO_STATIC_ROOT=/home/docs/kangaroo-review
Environment=KANGAROO_MOUNT=/kangaroo-review
Environment=KANGAROO_QUIET=1
ExecStart=/usr/bin/python3 /home/docs/kangaroo-review/server/metrics_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
